# marzban-inbound-supervisor.conf
# Supervisor config for auto-enabling Marzban inbounds
# Запускает скрипт через 60 секунд после старта Marzban

[program:marzban-auto-inbound]
command=/bin/bash -c "sleep 60 && python /app/marzban-auto-inbound.py"
autostart=true
autorestart=false
stderr_logfile=/var/log/marzban-auto-inbound.err.log
stdout_logfile=/var/log/marzban-auto-inbound.out.log
priority=999
